# conf/supervisord.conf
[program:gopkgs]
directory=/etc/gopkgs
command=/usr/bin/gopkgs serve

[program:gopkgs-queue]
directory=/etc/gopkgs
command=/usr/bin/gopkgs queue
